# supervisor/daemon.conf
[program:nginx]
command = /usr/sbin/nginx ;

[program:mongod]
command = mongod ;
